# docker/etc_master/supervisord.conf
[supervisord]
nodaemon=true

[program:salt-master]
command=/usr/bin/salt-master

[program:salt-api]
command=/usr/bin/salt-api